# deploy/ragsmith.service
# ─────────────────────────────────────────────────────────────────────────────
# RAGSmith – systemd service unit
#
# Install:
#   sudo cp deploy/ragsmith.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable ragsmith
#   sudo systemctl start ragsmith
#
# Logs:
#   sudo journalctl -u ragsmith -f
# ─────────────────────────────────────────────────────────────────────────────

[Unit]
Description=RAGSmith – Open-Source RAG Builder
After=network.target
Wants=network-online.target

[Service]
Type=exec
User=ragsmith
Group=ragsmith
WorkingDirectory=/home/ragsmith/app

# Load .env file — all secrets live here, NOT in this file
EnvironmentFile=/home/ragsmith/app/.env

# Gunicorn with uvicorn workers
ExecStart=/home/ragsmith/venv/bin/gunicorn main:app \
    --worker-class uvicorn.workers.UvicornWorker \
    --workers 1 \
    --bind 127.0.0.1:8000 \
    --timeout 120 \
    --access-logfile - \
    --error-logfile -

# Graceful reload on SIGHUP (used by deploy.sh)
ExecReload=/bin/kill -s HUP $MAINPID

# Restart policy
Restart=on-failure
RestartSec=5s
StartLimitBurst=5
StartLimitIntervalSec=60s

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ReadWritePaths=/home/ragsmith/app/data /home/ragsmith/app/exports

# Resource limits (appropriate for t2.micro: 1 vCPU, 1GB RAM)
LimitNOFILE=65536
MemoryMax=900M        # Leave ~100MB for OS
CPUQuota=95%

[Install]
WantedBy=multi-user.target
